# etc/init/iptables.conf
# Load saved iptables configuration

description "load iptables.conf"

start on startup

task
exec iptables-restore < /etc/iptables.conf
exec echo 'iptables configuration loaded!'
